# service/VenusOS-Tailscale/run
#!/bin/sh
exec 2>&1
exec /usr/bin/python3 /data/VenusOS-Tailscale/tailscale_control.py
