# service/VenusOS-Tailscale-backend/run
#!/bin/sh
exec /data/VenusOS-Tailscale/tailscaled \
  --no-logs-no-support \
  --netfilter-mode=off \
  --statedir=/data/conf/tailscale
